# plover-evdevd/plover-evdevd.service
[Unit]
Description=Plover evdev daemon

[Service]
ExecStart=/usr/bin/plover-evdevd
KillMode=process
Restart=always
